# scripts/interagent-daemon/com.unratified.interagent-webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unratified.interagent-webhook</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kashif/projects/unratified/scripts/interagent-daemon/webhook-listener.py</string>
        <string>--port</string>
        <string>8787</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/kashif/.nvm/versions/node/v24.14.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kashif/logs/interagent-webhook-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kashif/logs/interagent-webhook-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/kashif/projects/unratified</string>
</dict>
</plist>
