# scripts/interagent-daemon/interagent-webhook.service
[Unit]
Description=Interagent mesh webhook listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/kashif/projects/unratified/scripts/interagent-daemon/webhook-listener.py --port 8787
Environment=WEBHOOK_SECRET=%h/.config/interagent/webhook-secret
EnvironmentFile=-%h/.config/interagent/env
WorkingDirectory=/home/kashif/projects/unratified
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
